# pun.timer
[Unit]
Description=Runs pun daily

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
