# pun.service
[Unit]
Description=Pacman Update Notifer

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/pun

[Install]
WantedBy=pun.timer
